# com.lailabot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lailabot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ziyangliu/Github/zliu63/lailabot/.venv/bin/python</string>
        <string>-m</string>
        <string>lailabot</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ziyangliu/Github/zliu63/lailabot</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TELEGRAM_BOT_TOKEN</key>
        <string>YOUR_BOT_TOKEN_HERE</string>
        <key>TELEGRAM_USER_ID</key>
        <string>YOUR_USER_ID_HERE</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/ziyangliu/.lailabot/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ziyangliu/.lailabot/logs/stderr.log</string>
</dict>
</plist>
